# .config/crontab/.config
0 0/6 * * * make --silent -f ~/.config/dashy/Makefile cron
0 0 * * * make --silent -f ~/.config/myip/Makefile record-myip
0 9,11 * * * make --silent -f ~/.config/myip/Makefile check-if-ip-has-changed
0 0 * * * make --silent -f ~/.config/crontab/Makefile getconfig
