# .config/crontab/.config
0 0/6 * * * make --silent -f ~/.config/dashy/Makefile cron
