# .config/crontab/.config
0 0/6 * * * dashypublish >> ~/.config/dashy/crontab.logs 2>> ~/.config/dashy/crontab.logs
