# .config/crontab/.config
0 * * * * make --silent -f ~/projects/digital-garden/Makefile cron
0,30 9-13 * * * make --silent -f ~/.config/myip/Makefile check-if-ip-has-changed
0 0 * * * make --silent -f ~/Makefile cron-daily
0 16 * * * make --silent -f ~/.config/crontab/work.Makefile ts-fe-pull
