# .config/crontab/.config
0,30 9-13 * * * make --silent -f ~/.config/myip/Makefile check-if-ip-has-changed
0 0/6 * * * make --silent -f ~/.config/dashy/Makefile cron
0 0 * * * make --silent -f ~/.config/crontab/Makefile getconfig
0 0 * * * make --silent -f ~/projects/digital-garden cron
