# crontab/.config/crontab/.config
*/12 * * * * curl --location https://afk-slackbot.onrender.com
*/15 * * * * /opt/homebrew/bin/speedtest --secure --csv >> ~/.local/state/speedtest/speedtests.csv
0 * * * * make --silent -f ~/projects/digital-garden/Makefile cron
0 0 * * * make --silent -f ~/Makefile cron-daily
0 0 * * 0 make --silent -f ~/Makefile cron-weekly
@reboot sleep 15; make --silent -f ~/Makefile startup
